# systemd/docker-volume-iso.socket
[Unit]
Description=Docker ISO Volume plugin socket

[Socket]
ListenStream=/run/docker/plugins/iso.sock

[Install]
WantedBy=sockets.target
